# netshare-wan-watchdog.timer
[Unit]
Description=Verifica a cada 5 min se a WAN da NetShare está na banda 5 GHz

[Timer]
# Primeira verificação 2 min após o boot (dá tempo ao NM/Wi-Fi de subir)
OnBootSec=2min
# Depois a cada 5 min
OnUnitActiveSec=5min
Unit=netshare-wan-watchdog.service
# Persiste o calendário se o sistema esteve desligado
Persistent=true

[Install]
WantedBy=timers.target
