# netshare-mgmt-route.service
[Unit]
Description=NetShare management return route (keeps SSH/web alive when WAN is default)
After=NetworkManager.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/netshare/mgmt-route.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
